# package/debian/usr/lib/systemd/system/tibber.service
[Unit]
Description=Tibber Pulse app

[Service]
Type=simple
ExecStart=/usr/bin/tibber -c /etc/futurehome/tibber/config.json
User=tibber
WorkingDirectory=/var/lib/futurehome/tibber
Restart=always
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target